# home/fs/.config/autostart/startup_fs_http_selenium.desktop
[Desktop Entry]
Type=Application
Name=FS Startup Script
Exec=/home/fs/startup_fs_http_selenium.sh
X-GNOME-Autostart-enabled=true

